# ylx_files/x86/files/etc/uci-defaults/70-network-init
# #!/bin/sh

# uci add_list  network.@device[0].ports='eth1'
# uci set network.wan.device='eth2'
# uci set network.wan6.device='eth2'
# uci commit 
# /etc/init.d/network reload
