# package/ylx_CallManager/files/etc/uci-defaults/ylx_CallManagerDefaults
#!/bin/sh

echo "*/1 * * * * /sbin/ylx_ModemState_Detect.sh" | crontab -

exit 0